# service/bot/src/main/dist/bot.service
[Unit]
Requires = network-online.target
After = network-online.target

[Service]
Type = simple
ExecStart = /opt/cubing/bot/start-service.sh

[Install]
WantedBy = multi-user.target
